# root-ramips/usr/lib/opkg/info/hwnat.list
/etc/init.d/hwnat
/bin/hwnat
/bin/ac
/bin/acl
/bin/mtr
