# root-ramips/usr/lib/opkg/info/luci-mtk.list
/www/luci-static/resources/cbi/user.gif
/usr/lib/lua/luci/view/admin_network/wifi_status.htm
/usr/lib/lua/luci/view/admin_system/packages.htm
/usr/lib/lua/luci/store.lua
/usr/lib/lua/luci/model/cbi/admin_status/processes.lua
/usr/lib/lua/luci/view/cbi/dvalue.htm
/etc/config/luci
/usr/lib/lua/luci/view/cbi/upload.htm
/usr/lib/lua/luci/model/cbi/admin_system/leds.lua
/usr/lib/lua/bit.lua
/usr/lib/lua/luci/model/cbi/admin_system/ipkg.lua
/www/luci-static/resources/cbi/remove.gif
/www/luci-static/bootstrap/favicon.ico
/usr/lib/lua/luci/view/admin_status/dmesg.htm
/usr/lib/lua/luci/view/cbi/firewall_zoneforwards.htm
/usr/lib/lua/luci/view/admin_network/switch_status.htm
/www/luci-static/resources/icons/switch_disabled.png
/usr/lib/lua/luci/view/cbi/button.htm
/usr/lib/lua/luci/view/cbi/value.htm
/usr/lib/lua/luci/view/cbi/simpleform.htm
/www/luci-static/resources/wireless.svg
/usr/lib/lua/luci/i18n/base.en.lmo
/etc/init.d/luci_fixtime
/www/luci-static/bootstrap/html5.js
/usr/lib/lua/luci/template.lua
/usr/lib/lua/luci/model/ipkg.lua
/usr/lib/lua/luci/fs.lua
/www/luci-static/resources/icons/encryption_disabled.png
/usr/lib/lua/luci/view/cbi/map.htm
/www/luci-static/resources/icons/loading.gif
/usr/lib/lua/luci/view/header.htm
/usr/lib/lua/luci/view/cbi/delegator.htm
/www/luci-static/resources/icons/signal-75-100.png
/usr/lib/lua/nixio/util.lua
/usr/lib/lua/luci/model/cbi/admin_network/proto_static.lua
/www/luci-static/resources/icons/signal-none.png
/usr/lib/lua/luci/model/cbi/admin_system/system.lua
/usr/lib/lua/luci/view/footer.htm
/sbin/luci-reload
/usr/lib/lua/luci/model/cbi/admin_network/proto_ppp.lua
/usr/lib/lua/luci/view/admin_uci/changes.htm
/usr/lib/lua/luci/model/network.lua
/usr/lib/lua/luci/controller/admin/uci.lua
/www/luci-static/resources/icons/wifi.png
/www/luci-static/resources/cbi/save.gif
/usr/lib/lua/luci/view/error404.htm
/usr/lib/lua/luci/ip.lua
/usr/lib/lua/luci/debug.lua
/usr/lib/lua/luci/view/cbi/network_ifacelist.htm
/usr/lib/lua/luci/view/admin_system/flashops.htm
/usr/lib/lua/luci/view/admin_status/connections.htm
/usr/lib/lua/luci/sys/zoneinfo.lua
/www/luci-static/resources/icons/port_up.png
/usr/lib/lua/luci/config.lua
/usr/lib/lua/luci/view/cbi/lvalue.htm
/usr/lib/lua/luci/view/cbi/firewall_zonelist.htm
/usr/lib/lua/luci/util.lua
/usr/lib/lua/luci/sys.lua
/usr/lib/lua/luci/model/cbi/admin_system/backupfiles.lua
/www/luci-static/resources/cbi/file.gif
/usr/lib/lua/luci/view/admin_status/index.htm
/usr/lib/lua/nixio.so
/usr/lib/lua/luci/controller/admin/network.lua
/etc/config/ucitrack
/usr/lib/lua/luci/view/cbi/ucisection.htm
/usr/lib/lua/luci/view/cbi/full_valueheader.htm
/usr/lib/lua/luci/model/cbi/admin_system/fstab/mount.lua
/usr/lib/lua/luci/view/admin_system/upgrade.htm
/usr/lib/lua/luci/view/admin_system/ipkg.htm
/usr/lib/lua/luci/sys/zoneinfo/tzoffset.lua
/www/luci-static/resources/cbi/download.gif
/www/luci-static/resources/icons/port_down.png
/usr/lib/lua/luci/model/cbi/admin_network/routes.lua
/www/luci-static/resources/cbi/reset.gif
/usr/lib/lua/luci/view/cbi/full_valuefooter.htm
/www/luci-static/resources/cbi/help.gif
/usr/lib/lua/luci/sys/zoneinfo/tzdata.lua
/www/luci-static/resources/icons/tunnel.png
/www/luci-static/resources/xhr.js
/usr/lib/lua/luci/model/cbi/admin_network/network.lua
/www/luci-static/resources/icons/vlan_disabled.png
/usr/lib/lua/luci/http/protocol/date.lua
/usr/lib/lua/luci/cbi.lua
/usr/lib/lua/luci/view/admin_network/wifi_join.htm
/www/luci-static/resources/icons/vlan.png
/usr/lib/lua/nixio/fs.lua
/www/luci-static/resources/cbi/reload.gif
/usr/lib/lua/luci/view/admin_network/iface_overview.htm
/www/luci-static/resources/load.svg
/usr/lib/lua/luci/controller/admin/index.lua
/usr/lib/lua/luci/ltn12.lua
/usr/lib/lua/luci/model/cbi/admin_system/fstab/swap.lua
/usr/lib/lua/luci/model/cbi/admin_network/proto_pppoe.lua
/usr/lib/lua/luci/view/cbi/dynlist.htm
/usr/lib/lua/luci/view/admin_system/backupfiles.htm
/usr/lib/lua/luci/ccache.lua
/usr/lib/lua/luci/view/admin_system/clock_status.htm
/usr/lib/lua/luci/model/cbi/admin_system/startup.lua
/usr/lib/lua/luci/view/admin_status/syslog.htm
/usr/lib/lua/luci/controller/admin/system.lua
/usr/lib/lua/luci/model/cbi/admin_network/proto_pptp.lua
/usr/lib/lua/luci/view/cbi/tabmenu.htm
/usr/bin/luci-bwc
/www/luci-static/resources/cbi.js
/usr/lib/lua/luci/model/uci.lua
/lib/upgrade/luci-add-conffiles.sh
/www/luci-static/resources/cbi/find.gif
/usr/lib/lua/luci/sys/iptparser.lua
/usr/lib/lua/luci/cacheloader.lua
/www/luci-static/resources/bandwidth.svg
/usr/lib/lua/luci/view/cbi/header.htm
/www/cgi-bin/luci
/usr/lib/lua/luci/model/cbi/admin_network/iface_add.lua
/usr/lib/lua/luci/model/cbi/admin_network/proto_none.lua
/usr/lib/lua/luci/tools/webadmin.lua
/www/luci-static/resources/icons/signal-50-75.png
/usr/lib/lua/luci/controller/admin/status.lua
/usr/lib/lua/luci/view/cbi/cell_valueheader.htm
/usr/lib/lua/luci/init.lua
/usr/lib/lua/luci/view/cbi/footer.htm
/usr/lib/lua/luci/controller/admin/servicectl.lua
/usr/lib/lua/luci/view/cbi/tvalue.htm
/usr/lib/lua/luci/view/admin_status/iptables.htm
/www/luci-static/resources/cbi/edit.gif
/usr/lib/lua/luci/view/admin_status/routes.htm
/usr/lib/lua/luci/view/cbi/mvalue.htm
/usr/lib/lua/luci/view/cbi/cell_valuefooter.htm
/usr/lib/lua/luci/i18n.lua
/usr/lib/lua/luci/model/cbi/admin_system/admin.lua
/usr/lib/lua/luci/tools/status.lua
/usr/lib/lua/luci/view/cbi/tblsection.htm
/usr/lib/lua/luci/version.lua
/usr/lib/lua/luci/view/cbi/fvalue.htm
/usr/lib/lua/luci/model/cbi/admin_system/crontab.lua
/usr/lib/lua/luci/view/cbi/network_netinfo.htm
/etc/uci-defaults/luci-theme-bootstrap
/usr/lib/lua/luci/view/admin_system/applyreboot.htm
/usr/lib/lua/luci/view/sysauth.htm
/www/luci-static/resources/icons/wifi_disabled.png
/usr/lib/lua/luci/dispatcher.lua
/www/luci-static/resources/icons/signal-0-25.png
/usr/lib/lua/luci/model/cbi/admin_network/proto_dhcp.lua
/www/luci-static/resources/cbi/fieldadd.gif
/www/luci-static/bootstrap/mobile.css
/www/luci-static/resources/wifirate.svg
/usr/lib/lua/luci/view/themes/bootstrap/header.htm
/usr/lib/lua/luci/http.lua
/www/luci-static/resources/icons/wifi_big.png
/usr/lib/lua/luci/view/admin_status/bandwidth.htm
/usr/lib/lua/luci/sauth.lua
/usr/lib/lua/luci/model/cbi/admin_system/fstab.lua
/www/luci-static/resources/cbi/apply.gif
/www/luci-static/resources/cbi/key.gif
/usr/lib/lua/luci/view/admin_uci/apply.htm
/www/luci-static/resources/cbi/link.gif
/usr/lib/lua/luci/view/admin_status/wireless.htm
/usr/lib/lua/luci/view/cbi/nullsection.htm
/usr/lib/lua/luci/view/themes/bootstrap/footer.htm
/usr/lib/lua/luci/view/admin_network/diagnostics.htm
/usr/lib/lua/luci/model/cbi/admin_network/wifi.lua
/www/luci-static/resources/cbi/arrow.gif
/usr/lib/lua/luci/view/cbi/filebrowser.htm
/usr/lib/lua/luci/http/protocol/conditionals.lua
/usr/lib/lua/luci/http/protocol.lua
/usr/lib/lua/luci/view/cbi/tabcontainer.htm
/usr/lib/lua/luci/view/cbi/nsection.htm
/usr/lib/lua/luci/tools/proto.lua
/usr/lib/lua/luci/model/firewall.lua
/usr/lib/lua/luci/model/cbi/admin_network/dhcp.lua
/usr/lib/lua/luci/view/admin_status/load.htm
/www/luci-static/resources/icons/tunnel_disabled.png
/usr/lib/lua/luci/view/admin_system/reboot.htm
/usr/lib/lua/luci/http/protocol/mime.lua
/usr/lib/lua/luci/template/parser.so
/www/luci-static/resources/icons/signal-0.png
/usr/lib/lua/luci/view/cbi/apply_xhr.htm
/etc/init.d/luci_dhcp_migrate
/www/luci-static/resources/icons/switch.png
/usr/lib/lua/luci/model/cbi/admin_network/wifi_add.lua
/usr/lib/lua/luci/view/cbi/tsection.htm
/www/luci-static/resources/icons/encryption.png
/usr/lib/lua/luci/model/cbi/admin_network/hosts.lua
/usr/lib/lua/luci/controller/admin/filebrowser.lua
/www/luci-static/resources/icons/ethernet_disabled.png
/www/luci-static/resources/cbi/down.gif
/www/luci-static/resources/cbi/folder.gif
/www/luci-static/resources/icons/bridge_disabled.png
/usr/lib/lua/luci/cbi/datatypes.lua
/www/luci-static/resources/icons/wifi_big_disabled.png
/www/luci-static/resources/cbi/up.gif
/usr/lib/lua/luci/model/cbi/admin_network/vlan.lua
/usr/lib/lua/luci/view/error500.htm
/usr/lib/lua/luci/model/cbi/admin_network/proto_l2tp.lua
/usr/lib/lua/luci/view/admin_network/wifi_overview.htm
/www/luci-static/resources/icons/ethernet.png
/www/index.html
/www/luci-static/resources/icons/signal-25-50.png
/usr/lib/lua/luci/model/cbi/admin_network/proto_pppoa.lua
/usr/lib/lua/luci/view/cbi/browser.htm
/usr/lib/lua/luci/view/cbi/valueheader.htm
/usr/lib/lua/luci/view/admin_network/lease_status.htm
/www/luci-static/bootstrap/cascade.css
/www/luci-static/resources/connections.svg
/etc/uci-defaults/luci-i18n-english
/usr/lib/lua/luci/model/cbi/admin_system/buttons.lua
/usr/lib/lua/luci/model/cbi/admin_network/ifaces.lua
/usr/lib/lua/luci/view/admin_network/iface_status.htm
/usr/lib/lua/luci/view/cbi/compound.htm
/usr/lib/lua/luci/view/cbi/valuefooter.htm
/usr/lib/lua/luci/model/network/proto_ppp.lua
/usr/lib/lua/luci/view/indexer.htm
/usr/lib/lua/luci/view/admin_uci/revert.htm
/usr/lib/lua/luci/sgi/cgi.lua
/www/luci-static/resources/cbi/add.gif
/usr/lib/lua/luci/model/cbi/admin_network/proto_ahcp.lua
/www/luci-static/resources/icons/bridge.png
/usr/lib/lua/luci/view/admin_uci/changelog.htm
/usr/lib/lua/luci/view/cbi/network_netlist.htm
